# systemd/failing-disk-reporter.timer
[Unit]
Description=Failing Disk Reporter

[Timer]
OnCalendar=*-*-* 02:30:00
Unit=failing-disk-reporter.service

[Install]
WantedBy=timers.target
